# crontab
* * * * * python3 /Your_full_path/background.py
